# etc/systemd/system/fairyphone.service
# systemctl enable fairyphone.service
[Unit]
Description=Start your fairy phone
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart= /opt/phone/software/phone.sh
WorkingDirectory=/opt/phone/software
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=1s
[Install]
WantedBy=multi-user.target

